# script/mAlphaR22.cron
* */1 * * * /star/u/zchang/run22/cnipol/script/mAlphaR22.csh > /dev/null
