# script/mAsymR22.cron
* */1 * * * /star/u/zchang/run22/cnipol/script/mR22.csh > /dev/null
