# script/runAsymR22.cron
*/1 * * * * /star/u/xchu/run22/cnipol/script/r22asym.csh > /dev/null
