# script/runBlue.cron
* * * * * /star/u/xchu/run22/cnipol/script/mR22.csh > /dev/null
