# script/runYell.cron
*/1 * * * * /star/u/xchu/run22/cnipol/script/r22yell.csh > /dev/null
